# scanner/service/barcode-listener.service
[Unit]
Description=Barcode Scanner Listener Service
After=network.target

[Service]
Type=simple
User=root  # root権限が必要（入力デバイスへのアクセスのため）
WorkingDirectory=/path/to/home-inventory-manager/scanner
ExecStart=/usr/bin/python3 barcode_listener.py
Restart=on-failure
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=barcode-listener
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target 